# systemd/Morfeas_system.service.d/Morfeas_system.conf
[Service]
#Add absolute path of Morfeas_deamon configuration (.xml) file 
Environment="Morfeas_deamon_config_path="
