# systemd/Morfeas_system.service.d/Morfeas_system.conf
[Service]
Environment="Morfeas_deamon_config_path=/home/pi/Morfeas_proj/Morfeas_core/configuration/Morfeas_config.xml"
